# crontab
## Pseudo DDNS
*/20 * * * * root    /usr/local/bin/warning
